# deploy/systemd/lunch-order-monitor.service
[Unit]
Description=Lunch order health monitor + Telegram alert
After=docker.service lunch-order-stack.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/lunch_order/scripts/lunch_order_monitor.sh
